# config/includes.chroot/etc/rc.local
#!/bin/bash
/usr/local/bin/setup-wifi.sh
/usr/local/bin/setup-firewall.sh
exit 0

