# config/includes.chroot/etc/rc.local
#!/bin/bash
/usr/local/bin/setup-wifi.sh
exit 0

