# config/includes.chroot/etc/skel/.config/autostart/disable-icons.desktop
[Desktop Entry]
Type=Application
Name=Disable XFCE Desktop Icons
Exec=/home/student/.config/disable-icons.sh
X-GNOME-Autostart-enabled=true

