# config/includes.chroot/etc/xdg/autostart/disable-screensaver.desktop
[Desktop Entry]
Type=Application
Name=Disable Screensaver
Exec=xfconf-query -c xfce4-session -p /general/LockCommand -s ""

