# config/includes.chroot/etc/xdg/autostart/mitmproxy.desktop
[Desktop Entry]
Type=Application
Exec=mitmdump -p 8081 -s /usr/local/bin/whitelist.py
Hidden=false
NoDisplay=false
Name=MITMProxy Script
