# config/includes.chroot/etc/xdg/autostart/set-wallpaper.desktop
[Desktop Entry]
Type=Application
Exec=/usr/local/bin/set-wallpaper.sh
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=Set Wallpaper
Comment=Set wallpaper at startup
